# templates/armareforger-restart.service.j2
[Unit]
Description=Restart Arma Reforger Dedicated Server

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart armareforger.service
