# templates/armactl-bot.service.j2
[Unit]
Description=armactl Telegram Bot ({{ instance }})
After=network-online.target
Wants=network-online.target
RequiresMountsFor={{ project_root }} {{ bot_dir }}
StartLimitIntervalSec=0

[Service]
Type=simple
User={{ user }}
WorkingDirectory={{ project_root }}
ExecStart={{ python_bin }} -m armactl.telegram_bot --instance {{ instance }}
Restart=always
RestartSec=15
Environment=HOME={{ home_dir }}
Environment=USER={{ user }}
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
